# recipes-connectivity/wifi/files/maaxboard-wifi.service
[Unit]
Description=WIFI service

[Service]
Type=simple
ExecStart=/etc/wifi/wifi-start 
#Restart=on-failure

[Install]
WantedBy=multi-user.target
