# recipes-misc/wpa-conf/files/wpa-conf.timer
[Unit]
Description=Timer for the wpa-conf service

[Timer]
OnBootSec=20

[Install]
WantedBy=timers.target
