# recipes-misc/sshswitch/files/sshswitch.timer
[Unit]
Description=Timer for the sshswitch service

[Timer]
OnBootSec=20

[Install]
WantedBy=timers.target
