# recipes-misc/weston-timer/files/weston.timer
[Unit]
Description=Weston need start after system bootup

[Timer]
OnBootSec=1

[Install]
WantedBy=timers.target
